# system_files/mrcapuccin_gunicorn.socket
[Unit]
Description=chatgpt gunicorn socket

[Socket]
ListenStream=/run/chatgpt_gunicorn.sock

[Install]
WantedBy=sockets.target